# txqueuelen.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/txqueuelen.sh

[Install]
WantedBy=default.target
